# bin/com.mark.project-map.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mark.project-map</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/mark/bin/project-map-launch</string>
    </array>

    <!-- Run once daily at 07:00 (one-shot, not KeepAlive — no ThrottleInterval needed) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/project-map.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/project-map.err</string>
</dict>
</plist>
