# bin/com.mark.file-metadata-index.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mark.file-metadata-index</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/mark/bin/run-file-metadata-index-launch</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/mark/.claude/logs/file-metadata-index.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mark/.claude/logs/file-metadata-index.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/mark</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/mark</string>
        <key>PATH</key>
        <string>/Users/mark/.local/bin:/Library/Frameworks/Python.framework/Versions/3.13/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
